# systemd/triggerhappy.socket
[Socket]
ListenStream=/run/thd.socket

[Install]
WantedBy=sockets.target
