# udev/triggerhappy-udev.rules
# call the triggerhappy script if a new device is plugged in
ACTION=="add", SUBSYSTEM=="input", RUN+="/usr/sbin/th-cmd /var/run/thd.socket udev"
